# docs/launchd/com.polybot.paper-btc-15m.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.polybot.paper-btc-15m</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jeffrey/Documents/PolyBot/scripts/paper_btc_15m_launch.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jeffrey/Documents/PolyBot</string>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
